# tweet-pipeline/audit-summary.timer
[Unit]
Description=Run tweet audit summary every Sunday at 09:00 UTC

[Timer]
OnCalendar=Sun *-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
