# tweet-pipeline/scan-competitors.timer
[Unit]
Description=Scan Reddit for competitor mentions every 3 days at 10:00 UTC

[Timer]
OnCalendar=*-*-01/3 10:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
