# tweet-pipeline/scan-competitors.timer
[Unit]
Description=Scan Reddit for competitor mentions every 3 days at 10:00 UTC

[Timer]
# Run every 3 days: first run 5 min after enable, then every 3 days from last activation
OnActiveSec=5min
OnUnitActiveSec=3d
Persistent=true

[Install]
WantedBy=timers.target
