# tweet-pipeline/scan-newsletter.service
[Unit]
Description=Newsletter scanner — score stories for tweet/blog pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/tweet-pipeline/scan-newsletter.sh
WorkingDirectory=/opt/tweet-pipeline
EnvironmentFile=/opt/blog-pipeline.env
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
