# tweet-pipeline/tweet-engagement.service
[Unit]
Description=Twitter engagement metrics bridge to PostHog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/ethernal-blog-stack/tweet-pipeline/engagement-bridge.sh
WorkingDirectory=/opt/ethernal-blog-stack
EnvironmentFile=/opt/blog-pipeline.env
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
