# tweet-pipeline/tweet-publish.service
[Unit]
Description=Tweet publisher — post queued tweets past their scheduled time
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/tweet-pipeline/publish.sh
WorkingDirectory=/opt/tweet-pipeline
EnvironmentFile=/opt/blog-pipeline.env
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
