# tweet-pipeline/tweet-draft.service
[Unit]
Description=Tweet draft pipeline — pick source, research, draft, queue
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/tweet-pipeline/draft.sh
WorkingDirectory=/opt/tweet-pipeline
EnvironmentFile=/opt/blog-pipeline.env
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
