# tweet-pipeline/tweet-engagement.service
[Unit]
Description=Twitter engagement metrics bridge to PostHog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=blog
Group=blog
ExecStart=/opt/tweet-pipeline/engagement-bridge.sh
WorkingDirectory=/opt/tweet-pipeline
EnvironmentFile=/opt/blog-pipeline.env
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
